# templates/default/activemq.service.erb
[Unit]
Description=ActiveMQ service
After=network.target

[Service]
Type=forking
ExecStart=<%= @activemq_base %>/bin/activemq start
ExecStop=<%= @activemq_base %>/bin/activemq stop
User=root
Group=root
Restart=always
RestartSec=9
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=activemq

[Install]
WantedBy=multi-user.target
