# templates/default/etc-init.d-opencast.erb
#!/bin/bash
#
# opencast     Start/Stop Opencast
#
# chkconfig:   - 60 40
# description: Opencast is a free, open-source platform to support the \
#              management of educational audio and video content. Institutions \
#              will use Opencast to produce lecture recordings, manage \
#              existing video, serve designated distribution channels, and \
#              provide user interfaces to engage students with educational \
#              videos.

### BEGIN INIT INFO
# Provides: opencast
# Required-Start: $local_fs $remote_fs $syslog $network
# Required-Stop:
# Default-Start:
# Default-Stop:
# Short-Description: run opencast
# Description: Opencast is a free, open-source platform to support the
#              management of educational audio and video content. Institutions
#              will use Opencast to produce lecture recordings, manage
#              existing video, serve designated distribution channels, and
#              provide user interfaces to engage students with educational
#              videos.
### END INIT INFO

opencast="<%= @opencast_root %>"
prog="opencast"
user="opencast"
lockfile=/var/lock/subsys/${prog}
[ -d "/var/lock/subsys" ] || lockfile="/var/lock/LCK.${prog}"

# Used in the "restart until the daemon is actually working" feature
# now included in start()
test_string="j_username"
max_sleep_attempts=10
max_start_attempts=3
attempts=0
sleep_attempts=0
sleeptime=10

# Load configuration files
[ -e /etc/sysconfig/$prog ] && . /etc/sysconfig/$prog

success() {
   printf "\r%-58s [\033[32m  OK  \033[0m]\n" "$1"
}

failed() {
   printf "\r%-58s [\033[31mFAILED\033[0m]\n" "$1"
}

start() {

   # If you don't want to boot opencast automatically after instance boot, set
   # "dont_start_opencast_automatically" to "true" in your stack's custom_json.
   # *THIS WILL KEEP YOU FROM BOOTING OPENCAST ENTIRELY UNTIL YOU UNSET THE VARIABLE
   # AND UPDATE YOUR CHEF RECIPES TO UPDATE YOUR JSON METADATA*. Use with caution.
   # This setting can be overridden on a per-instance basis by uncommenting the
   # `FORCE_OPENCAST_START=1` line in /etc/sysconfig/opencast
   #
   if [ -z "$FORCE_OPENCAST_START" ] && [ $(/usr/sbin/opsworks-agent-cli get_json | /usr/bin/jq '.dont_start_opencast_automatically' | grep 'true') ];
   then
     echo "Not starting because dont_start_opencast_automatically is set in the stack's custom_json"
     exit 0;
   fi

   smsg="Starting $prog: "
   echo -n "$smsg"

   # Start Opencast and create a lockfile
   ( flock -n 9 && sudo -u $user $opencast/bin/start-opencast server & > /dev/null ) 9> $lockfile
   retval=$?

   # If we failed with retval=1, Opencast might be already up and running. In
   # that case we still want to return a success:
   [ $retval -eq 1 ] && rh_status_q && echo && exit 1

   # If we failed to start Opencast but a lock was created, we want to remove
   # the file (flock does not remove the file automatically):
   [ ! $retval -eq 0 ] && rm -f $lockfile

   # Attempt to restart the daemon until it actually looks like it's
   # started correctly
   started=0;
   # sleep attempts first
   while [ "$sleep_attempts" -lt "$max_sleep_attempts" ] && [ "$started" -eq 0 ]; do
     sleep $sleeptime
     foundLoginPage= $(/usr/bin/curl -s -L http://localhost/ | grep -q "$test_string");
     if ! $foundLoginPage; then
       echo "$(date) could not find $test_string at localhost, waiting another $sleeptime seconds "
       sleep_attempts=$[$sleep_attempts + 1]
     else
       echo "$(date) Found $test_string at localhost"
       started=1;
     fi
   done
   # if sleeping didn't work, try a restart
   if [ "$attempts" -lt "$max_start_attempts" ] && [ "$started" -eq 0 ] ; then
     sleep_attempts=0
     attempts=$[$attempts + 1]
     smsg="Daemon did not start after $sleep_attempts sleep attempts, trying a restart"
     restart
   fi

   if [ "$attempts" -ge "$max_start_attempts" ] && [ "$started" -eq 0 ] ; then
     retval=1
   fi

   [ $retval -eq 0 ] && success "$smsg" || failed "$smsg"
   return $retval

}

stop() {
   smsg="Stopping $prog: "
   echo -n "$smsg"
   [ -f $lockfile ] && $opencast/bin/stop-opencast
   retval=$?
   if [ $retval -eq 0 ]
   then
      rm -f $lockfile
      success "$smsg"
   else
      failed "$smsg"
   fi
   return $retval
}

restart() {
   stop
   start
}

force_reload() {
   restart
}

rh_status() {
   # run checks to determine if the service is running or use generic status
   if [ -f $lockfile ]
   then
      sudo -u $user $opencast/bin/start-opencast status
      retval=$?
      if [ $retval -eq 0 ]
      then
         return 0
      else
         echo "${prog} dead but subsys locked" && return 2
      fi
  fi
  echo "${prog} is stopped"
  return 3
}

rh_status_q() {
   rh_status >/dev/null 2>&1
}


case "$1" in
   start)
      rh_status_q && exit 0
      $1
      ;;
   stop)
      rh_status_q || exit 0
      $1
      ;;
   restart)
      $1
      ;;
   reload)
      rh_status_q || exit 7
      ;;
   force-reload)
      force_reload
      ;;
   status)
      rh_status
      ;;
   condrestart|try-restart)
      rh_status_q || exit 0
      restart
      ;;
   *)
      echo $"Usage: $0 {start|stop|status|restart|condrestart|try-restart|reload|force-reload}"
      exit 2
esac
exit $?
